# games/FarmingSimulator25/supervisor/kasmvnc.conf
[program:kasmvnc]
command=/usr/bin/kasmvncserver --port --{{VNC_PORT}}
autostart=true
autorestart=true
stdout_logfile=/var/log/kasmvnc.log
stderr_logfile=/var/log/kasmvnc_error.log
